# server/know-gateway.service.example
# Example systemd unit for running know on a plain Linux host (not a Sprite).
# Copy to /etc/systemd/system/know-gateway.service, edit paths/env, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now know-gateway
#
# Schedule the heartbeat (off-box reconcile) hourly via cron. Carry the env vars
# wake.py needs (scheduler runs outside the service env), e.g.:
#   0 * * * * KNOW_NAME=know KNOW_MODEL=claude-haiku-4-5-20251001 KNOW_ALERT_WEBHOOK=... /home/know/know-gateway/.venv/bin/python /home/know/know-gateway/wake.py

[Unit]
Description=know team brain gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=know
Group=know
WorkingDirectory=/home/know/know-gateway
# No API key here — the service inherits the box's Claude auth (whatever makes
# interactive `claude` work for this user: subscription token, API key, or ~/.claude).
Environment=KNOW_SECRET=replace-with-your-mcp-secret
Environment=KNOW_NAME=know
Environment=KNOW_MODEL=claude-haiku-4-5-20251001
ExecStart=/home/know/know-gateway/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
